# dot_config/exact_systemd/exact_user/exact_wayland-wm@.service.d/override.conf
[Service]
# Protect from systemd-oomd (Fedora, Ubuntu, Arch with systemd-oomd)
OOMPolicy=continue

# Protect from Kernel OOM Killer
# Note: As a user service, you normally cannot set this to -1000 
# without root delegation, but -250 to -500 usually works best.
OOMScoreAdjust=-500
