# dot_config/exact_systemd/exact_user/exact_wayland-wm@hyprland.desktop.service.d/override.conf
[Service]
# -1000 means "never kill this process"
# -500 means "kill this only if absolutely necessary"
OOMScoreAdjust=-1000
